# images/loadbalancer/update-config.timer
[Unit]
Description=update nginx config

[Timer]
OnUnitActiveSec=5s
OnBootSec=5s

[Install]
WantedBy=timers.target
