# images/tags-web-server/tags-web-server.service
[Unit]
After=network.service

[Service]
User=ubuntu
ExecStart=/usr/local/bin/tags-web-server.sh

[Install]
WantedBy=default.target
